# supervisor/php-fpm.conf
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
user=root
stdout_logfile=/var/www/html/storage/logs/php-fpm.log
stderr_logfile=/var/www/html/storage/logs/php-fpm.err.log